# packaging/debian/postrm
#!/bin/sh
set -e

case "$1" in
  purge)
    if command -v systemctl >/dev/null 2>&1; then
      systemctl disable zoeken.service >/dev/null 2>&1 || true
      systemctl daemon-reload || true
    fi
    # Leave /etc/zoeken and /var/lib/zoeken for the operator on remove;
    # only purge removes the system user/group when empty of conffiles.
    ;;
  remove)
    if command -v systemctl >/dev/null 2>&1; then
      systemctl daemon-reload || true
    fi
    ;;
esac

exit 0
